# docker-prod/supervisor/conf/laravel.conf
[program:laravel]
process_name=%(program_name)s_%(process_num)02d
command=php artisan serve --host=127.0.0.1 --port=9000
directory=/var/www/api
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/laravel.log
